# musicserver/s6-rc.d/02-avahi/run
#!/command/with-contenv sh
echo "Starting avahi"
exec avahi-daemon --no-chroot
